# nanonote-files/script-files/etc/uci-defaults/01-optimization-for-ben-nanonote
#!/bin/sh

/etc/init.d/atd disable
/etc/init.d/led disable
/etc/init.d/telnet disable
