# keymouse/files/keymouse.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
START=43

start() {
	start-stop-daemon -b -S -x /usr/bin/keymouse -- /etc/keymouse.cfg 0
}
