# nanonote-files/script-files/etc/init.d/ben-nanonote
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
	/usr/sbin/alsactl restore -f /etc/asound.state
}

stop() {
	/usr/sbin/alsactl store -f /etc/asound.state
}
