# fn-handle/files/fn-handle.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
START=40

DEVICE_PW="/dev/input/event0"
DEVICE_KB="/dev/input/event1"

start() {
	fnkeys $DEVICE_PW &
	fnkeys $DEVICE_KB &
}

stop() {
	killall fnkeys
}
