# fn-handle/files/fn-handle.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
START=40

DEVICE_PW="/dev/input/event0"

start() {
	start-stop-daemon -S -x /usr/bin/fnkeys -b $DEVICE_PW
}

stop() {
	killall fnkeys
}
