# fn-handle/files/fn-handle.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org
START=40

start() {
	fn-handle
}

stop() {
	killall fn-handle
}
