# nanonote-files/data/qi_lb60/files/etc/uci-defaults/02-gforth-work-around
#!/bin/sh

(cd /usr/lib; ln -s libffi.so.5.0.10 libffi.so.4)
