# nanonote-files/data/qi_lb60/files/etc/uci-defaults/03-python-work-around
#!/bin/sh

rm -f /usr/lib/python2.6/lib-dynload/readline.so
